# cluster/saltbase/salt/flanneld/mk-docker-opts.conf.upstart
description "Generate DOCKER_OPTS uses Flannel network"

respawn

start on flanneld

script
  while [ ! -e "/run/flannel/subnet.env" ]
  do
    sleep 3
  done

  /usr/bin/mk-docker-opts.sh -d /run/flannel/docker
end script

